# cloudflare-ddns.timer
[Unit]
Description=Check my ip address

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=cloudflare-ddns.service

[Install]
WantedBy=timers.targe
