# systemd/cloudflare-ddns-aaaa@.timer
[Unit]
Description=Check my ip address

[Timer]
OnBootSec=1m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
