# systemd/cloudflare-ddns-a@.timer
[Unit]
Description=Check my ip address

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
